# failover-monitor.service
[Unit]
Description=Cloudflare Tunnel Monitor
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /opt/failover-monitor/failover_monitor.py
Restart=always
RestartSec=10
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
